# services/mpd/finish
#!/bin/bash
echo 'mpd exited... restarting'
