# services/mpd/run
#!/usr/bin/with-contenv /bin/bash
echo 'Starting mpd'
exec mpd --no-daemon --verbose
