# services/tail/run
#!/usr/bin/with-contenv /bin/bash

exec tail -F /tmp/mpd.log
